# surface-acpi-watcher.service
[Unit]
Description=Surface ACPI Quell — watcher that checks & repairs the fix
Documentation=file:/usr/local/share/doc/surface-acpi-quell/README.md
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/surface-acpi-watcher
StandardOutput=journal
StandardError=journal
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
CapabilityBoundingSet=CAP_SYS_MODULE
# ProtectSystem=strict makes the whole tree read-only; the watcher still needs
# a writable /tmp (auto_detect_gpes mktemp) and its data dir (history/state).
PrivateTmp=yes
StateDirectory=surface-acpi-quell
ReadWritePaths=/var/lib/surface-acpi-quell

[Install]
WantedBy=multi-user.target
